# deploy/aiops-python-sidecar.service
[Unit]
Description=AIOps Python AI Sidecar (LangGraph + Pipeline Executor, port 8050)
After=network.target aiops-java-api.service
# Sidecar calls back into Java for every DB op — wait for it to be up.
Requires=aiops-java-api.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/aiops/python_ai_sidecar
EnvironmentFile=/opt/aiops/python_ai_sidecar/.env
ExecStart=/opt/aiops/venv_sidecar/bin/uvicorn python_ai_sidecar.main:app \
  --host 127.0.0.1 --port 8050 --workers 1
ExecStop=/bin/sh -c 'lsof -ti tcp:8050 | xargs kill -TERM 2>/dev/null || true'
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiops-python-sidecar

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/aiops /tmp /var/log/aiops

[Install]
WantedBy=multi-user.target
